# dist/mastodon-metrics-collector.service
[Unit]
Description=mastodon-metrics
After=network.target

[Service]
Type=simple
User=mastodon
WorkingDirectory=/home/mastodon/live
Environment="RAILS_ENV=production"
ExecStart=/home/mastodon/.rbenv/shims/bundle exec bin/metrics-collector
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
