# dist/mastodon-streaming.service
[Unit]
Description=mastodon-streaming
After=network.target
Wants=mastodon-streaming@4000.service

[Service]
Type=oneshot
ExecStart=/bin/echo "mastodon-streaming exists only to collectively start and stop mastodon-streaming@ instances, shimming over the migration to templated mastodon-streaming systemd unit"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
